# usb_rest_server.service
[Unit]
Description=USB Media Remote
After=multi-user.target

[Service]
ExecStart=/usr/bin/UsbHidRest.py

[Install]
WantedBy=multi-user.target
