# usb_hid_server.service
[Unit]
Description=USB Media Remote
After=multi-user.target

[Service]
ExecStart=/usr/bin/UsbHidServer.py

[Install]
WantedBy=multi-user.target
