# usb_rest_server.service
[Unit]
Description=USB Media Remote
After=multi-user.target

[Service]
ExecStart=/usr/local/bin/UsbHidRest

[Install]
WantedBy=multi-user.target
